# nginwho.service
[Unit]
Description=nginx log parser and analyzer
After=Network.target

[Service]
Type=simple
RemainAfterExit=yes
ExecStart=/usr/local/bin/nginwho --logPath:/var/log/nginx/access.log --dbPath:/var/log/nginwho.db --show-real-ips:true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
